# paramiko/django_deploy/polls.ini
[program:polls]
command=/root/pythonenv/django1.8/bin/python /projects/polls/tornado_server.py
directory=/projects/polls